# etc/supervisor/conf.d/supervisord.conf
[supervisord]
logfile = /var/log/supervisord/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = true

[program:elasticsearch-config-writer]
command=node /opt/elasticsearch-config-writer
priority=1
autorestart=false

[program:elasticsearch]
command=/usr/share/elasticsearch/bin
priority=999
stdout_logfile=/var/log/elasticsearch/stdout.log
stderr_logfile=/var/log/elasticsearch/stderr.log